# templates/default/nagios-nsca-server.erb
#!/bin/sh
#
# Autogenerated by Chef.
#
#  Created 2013 Miguel Armas <kuko@canarytek.com>
#  Based on script by jaclu@grm.se
#
# nagios-nsca-server This shell script takes care of starting and stopping
#               nsca.
#
# chkconfig: 2345 80 30
# description: nsca is a daemon for receiving asyncronous nagios notifications
# processname: nsca
# config: /etc/nagios/nsca.cfg


# Source function library
if [ -f /etc/rc.d/init.d/functions ]; then
. /etc/rc.d/init.d/functions
elif [ -f /etc/init.d/functions ]; then
. /etc/init.d/functions
elif [ -f /etc/rc.d/functions ]; then
. /etc/rc.d/functions
fi

<% if ['rhel', 'fedora'].include?(node['platform_family']) -%>
# Source networking configuration.
. /etc/sysconfig/network

# Check that networking is up.
[ ${NETWORKING} = "no" ] && exit 0
<% end -%>

NscaBin=/usr/sbin/nsca
NscaCfg=<%= node['nagios']['nsca_server']['conf_dir'] %>/nsca.cfg
LockFile=/var/lock/subsys/nsca

# See how we were called.
case "$1" in
  start)
  # Start daemons.
  echo -n "Starting nagios-nsca-server: "
  daemon $NscaBin -c $NscaCfg -d
  echo
  touch $LockFile
  ;;
  stop)
  # Stop daemons.
  echo -n "Shutting down nagios-nsca-server: "
  killproc nsca
  echo
  rm -f $LockFile
  ;;
  restart)
  $0 stop
  $0 start
  ;;
  status)
  status nsca
  ;;
  *)
  echo "Usage: nagios-nsca-server {start|stop|restart|status}"
  exit 1
esac

exit 0
